# docs/pyronms-network-poller.service
[Unit]
Description=PyroNMS Network Graphs Poller (Cacti-style, multi-vendor SNMP)
After=network.target influxdb.service
Wants=influxdb.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/PyroNMS-repo
ExecStart=/usr/bin/python3 /root/PyroNMS-repo/workers/network_poller.py
Environment=NET_POLLER_THREADS=8
Environment=NETWORK_INFLUX_BUCKET=network_monitoring
Restart=on-failure
RestartSec=15
StandardOutput=append:/opt/ont-monitor/logs/network_poller.log
StandardError=append:/opt/ont-monitor/logs/network_poller.log

[Install]
WantedBy=multi-user.target
